# announceflow.service
[Unit]
Description=AnnounceFlow - Scheduled Audio System
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/announceflow
ExecStart=/usr/bin/python3 /home/admin/announceflow/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
